# pi/stage/etc/systemd/system/kiln-controller.service
[Unit]
Description=Kiln Controller Service
After=network.target

[Service]
# Service runs as 'pi' user (standard on RPi)
User=pi
Group=dialout

# Working directory on the target machine
WorkingDirectory=/opt/makerspace-kiln/service

# Start command
ExecStart=/usr/bin/node index.js

# Restart automatically if it crashes
Restart=always
RestartSec=10

# Environment config
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target